# scripts/okr-backup.service
[Unit]
Description=OKR Database Backup Service
After=network.target postgresql.service

[Service]
Type=oneshot
User=okr-backup
Group=okr-backup
Environment="NODE_ENV=production"
EnvironmentFile=/etc/okr/.env
WorkingDirectory=/opt/okr
ExecStart=/usr/bin/node /opt/okr/scripts/backup-db.js
StandardOutput=journal
StandardError=journal

# Security options
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
SystemCallArchitectures=native
MemoryDenyWriteExecute=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
